# files/crontab
0 6 * * * sh /opt/update_infosec_csv.sh
