# files/crontab
0 * * * * sh /opt/update_infosec_csv.sh
